# home/zetxx/crontab
*/5 * * * * /home/zetxx/router-rpi-4G/router/modem/router-rpi-4g-status.sh
