# system/home/zetxx/crontab/load
*/5 * * * * /home/zetxx/router-rpi-4G/system/home/zetxx/crontab/router-rpi-4g-status.sh
